# crontabui/crontabui.service
[Unit]
Description=Crontab UI

[Service]
User=root
Type=simple
ExecStart=/opt/simplecloud/crontabui/run.sh
