# thelounge/thelounge.service
[Unit]
Description=TheLounge

[Service]
User=thelounge
Type=simple
ExecStart=/opt/simplecloud/thelounge/run.sh
