# matterbridge/matterbridge.service
[Unit]
Description=Matterbridge

[Service]
User=matterbridge
Type=simple
ExecStart=/opt/simplecloud/matterbridge/run.sh
