# radicale/radicale.service
[Unit]
Description=Radicale

[Service]
User=radicale
Type=simple
ExecStart=/opt/simplecloud/radicale/run.sh
