# peercalls/peercalls.service
[Unit]
Description=Peercalls

[Service]
User=peercalls
ExecStart=/opt/simplecloud/peercalls/run.sh
